# borg/crontab
0 */6 * * *     /usr/local/bin/create.sh
0 1 * * Sun     /usr/local/bin/prune.sh
